# documentation/supervisor/octane.ini
[program:octane-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan octane:start --server=swoole --workers=16 --port=8000
autostart=true
autorestart=true
user=daemon
redirect_stderr=true
stdout_logfile=var/log/laravel-octane-worker.log​
stopwaitsecs=3600
