# documentation/supervisor/queue.ini

[program:queue-worker-q1]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q1 --sleep=5 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-q1.log​

[program:queue-worker-q2]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q2 --sleep=5 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-q2.log​

[program:queue-worker-q3]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q3 --sleep=5 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-q3.log​
