# documentation/supervisor/queue.ini

[program:queue-worker-a]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qa --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-a.log​

[program:queue-worker-b]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qb --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-b.log​

[program:queue-worker-c]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qc --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-c.log​

[program:queue-worker-d]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qd --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-d.log​

[program:queue-worker-e]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qe --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-e.log​

[program:queue-worker-f]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qf --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-f.log​

[program:queue-worker-g]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qg --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-g.log​

[program:queue-worker-h]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qh --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-h.log​

[program:queue-worker-i]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qi --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-i.log​

[program:queue-worker-j]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qj --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-j.log​

[program:queue-worker-k]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qk --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-k.log​

[program:queue-worker-l]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=ql --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-l.log​

[program:queue-worker-m]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qm --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-m.log​

[program:queue-worker-n]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qn --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-n.log​

[program:queue-worker-o]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qo --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-o.log​

[program:queue-worker-p]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qp --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-p.log​

[program:queue-worker-q]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qq --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-q.log​

[program:queue-worker-r]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qr --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-r.log​

[program:queue-worker-s]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qs --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-s.log​

[program:queue-worker-t]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qt --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-t.log​

[program:queue-worker-u]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qu --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-u.log​

[program:queue-worker-v]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qv --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-v.log​

[program:queue-worker-w]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qw --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-w.log​

[program:queue-worker-x]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qx --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-x.log​

[program:queue-worker-y]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qy --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-y.log​

[program:queue-worker-z]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=qz --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-z.log​

[program:queue-worker-0]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q0 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-0.log​

[program:queue-worker-1]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q1 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-1.log​

[program:queue-worker-2]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q2 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-2.log​

[program:queue-worker-3]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q3 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-3.log​

[program:queue-worker-4]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q4 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-4.log​

[program:queue-worker-5]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q5 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-5.log​

[program:queue-worker-6]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q6 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-6.log​

[program:queue-worker-7]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q7 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-7.log​

[program:queue-worker-8]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q8 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-8.log​

[program:queue-worker-9]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/nginx/htdocs/xwa/artisan queue:work --queue=q9 --sleep=3 --tries=1
autostart=true
autorestart=true
user=daemon
numprocs=1
redirect_stderr=true
stdout_logfile=var/log/supervisor/xwa-queue-worker-9.log​
